# supervisord.conf
; supervisord.conf — process manager for OrganistBot
;
; Usage:
;   supervisord -c supervisord.conf          # start in background
;   supervisorctl -c supervisord.conf status # check status
;   supervisorctl -c supervisord.conf stop all
;   supervisorctl -c supervisord.conf restart all
;
; %(here)s expands to the directory containing this file (the project root),
; so no absolute paths are needed.

[unix_http_server]
file=%(here)s/logs/supervisor.sock

[supervisord]
logfile=%(here)s/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=%(here)s/logs/supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/logs/supervisor.sock

[program:scraper]
command=%(here)s/scripts/run-scraper.sh
directory=%(here)s
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=%(here)s/logs/scraper.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=%(here)s/logs/scraper-err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3

[program:telegram-bot]
command=%(here)s/scripts/run-telegram-bot.sh
directory=%(here)s
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=%(here)s/logs/telegram.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=%(here)s/logs/telegram-err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
